# debian/radeon-profile-daemon.postrm
#!/bin/sh
set -x

echo "Enabling radeon profile daemon"
systemctl enable radeon-profile-daemon.service
wait $!

exit 0
