# radeon-profile-daemon/extra/radeon-profile-daemon.conf
description "radeon-profile-daemon Upstart"
author "gogo"

# start on startup
start on runlevel [2345]
stop on shutdown|reboot

exec radeon-profile-daemon

respawn


